# etc/init/mr808-sc.conf
#!upstart
description "MR-808 SClang Server"

env USER=root

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

respawn

exec start-stop-daemon --start --make-pidfile --pidfile /var/run/mr808-sc.pid --chuid $USER --exec /root/MR-808/sc/mr808.sh >> /var/log/mr808-sc.log 2>&1

