# etc/init/mr808-node.conf
#!upstart
description "MR-808 NodeJS Server"

env USER=root

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

respawn

exec start-stop-daemon --start --make-pidfile --pidfile /var/run/mr808-node.pid --chuid $USER --exec /root/MR-808/mr808.sh >> /var/log/mr808-nodejs.log 2>&1

